# devices/phytiumpi-hooks/phytiumpi-pwm/pwm-init.service
[Unit]
[Service]
Type=forking
ExecStart=/usr/bin/phytiumpi-pwm
User=root
[Install]
WantedBy=multi-user.target

